# .config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=ibus-daemon -drxR
Icon=
Name=ibus-daemon
Path=
Terminal=False
Type=Application
